# config/gentoo/etc/rc.local
#!/bin/sh
if [ "x$lyt" = "x" ];then
	echo "null" > /tmp/lyt
else
	echo $lyt > /tmp/lyt
fi
